# sage/geeko-sage.service
[Unit]
Description=Le Sage de Geeko (legacy fortune service)
After=network.target

[Service]
Type=simple
Environment=GEEKO_SAGE_FLAG_PATH=/etc/geeko-sage/patched
ExecStart=/opt/geeko-sage/venv/bin/uvicorn sage_service.main:app --host 0.0.0.0 --port 9000
WorkingDirectory=/opt/geeko-sage
Restart=on-failure

[Install]
WantedBy=multi-user.target
